# ensure-nix-dir.service
[Unit]
Description=Ensure /nix is present
[Service]
Type=oneshot
ExecStart=mkdir -p -m 0755 /nix
ExecStart=chown -R bjk /nix
